# web-messaging/crontab
*/5 * * * * /var/www/html/signal-parsing-webmessaging.sh
